# deploy/selfhosted/systemd/yetibot.service.d/20-gemini.conf
[Service]
Environment=GEMINI_CLI_TRUST_WORKSPACE=true
